# zca-bot/cron-daily-report.timer
[Unit]
Description=CRM Thaco Auto - báo cáo cuối ngày 18:00

[Timer]
OnCalendar=Mon..Sat 18:00:00
Persistent=false

[Install]
WantedBy=timers.target
