# zca-bot/cron-reminders.timer
[Unit]
Description=CRM Thaco Auto - reminders mỗi giờ trong giờ hành chính

[Timer]
OnCalendar=Mon..Sat 08..18:00:00 Asia/Ho_Chi_Minh
Persistent=false

[Install]
WantedBy=timers.target
